# backup/vps/fs/etc/cron.daily/molb
#!/bin/bash
su - franck -c "pg_dump molb | gzip > /home/franck/backup/db/molb-$(date +%Y-%m-%d).gz"
