# backup/vps/fs/etc/cron.daily/molb
#!/bin/bash
su - molb -c "pg_dump molb | gzip > /home/molb/backup/db/molb-$(date +%Y-%m-%d).gz"
